# deploy/systemd/aiorahub.service
[Unit]
Description=AI Ora Hub Spring Boot App
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=aiorahub
Group=aiorahub
WorkingDirectory=/opt/aiorahub
EnvironmentFile=/etc/aiorahub/aiorahub.env
ExecStart=/usr/bin/java -jar /opt/aiorahub/target/truehire-0.0.1-SNAPSHOT.jar
SuccessExitStatus=143
Restart=always
RestartSec=10
UMask=0027

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
IPAddressDeny=any
IPAddressAllow=localhost
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
